# step2/roles/redfish/templates/server.conf.j2
[program:server]
command=/usr/bin/python3.6 /opt/redfish/server.py
autostart=true
autorestart=true