# step3/roles/redfish/files/conf/server.conf
[program:server]
directory=/opt/redfish
command=/usr/bin/gunicorn-3.6 server:app -b localhost:8000
autostart=true
autorestart=true