# step4/roles/redfish/files/redfish.conf
[program:redfish]
command=/usr/bin/gunicorn-3.6 redfish:app -b unix:/opt/redfish/run/gunicorn.sock
directory=/opt/redfish/
logfile=/var/log/redfish/supervisor.log
autostart=true
autorestart=true